# templates/managed-user-pod.service.erb
[Unit]
Description=User Pod <%= @name %>
After=network.target

<% if @has_container_disk -%>
Requires=var-lib-containers.mount
After=var-lib-containers.mount
<% end
@require_mounts.each do |rm| -%>
Requires=<%= rm.split('/')[1..-1].map{|s| s.gsub('-','\x2d') }.join('-') %>.mount
After=<%= rm.split('/')[1..-1].map{|s| s.gsub('-','\x2d') }.join('-') %>.mount
<% end -%>

[Service]
Type=forking
KillMode=none
Restart=always
PIDFile=/run/pods/<%= @uid %>/<%= @unique_name %>.pid

User=<%= @user %>
Group=<%= @real_group %>

Environment="XDG_RUNTIME_DIR=/run/pods/<%= @uid %>" "REGISTRY_AUTH_FILE=/var/lib/containers/users/<%= @user %>/data/auth.json"

SyslogIdentifier=<%= @unique_name %>

ExecStartPre=/usr/local/bin/update-container-auth.sh <%= @user %>
<% unless @run_flags['skip_image_pull'] -%>
ExecStartPre=-/usr/local/bin/pod-update-image.sh <%= @pod_yaml_path %>
<% end -%>
ExecStartPre=-/usr/bin/podman pod rm -f <%= @sanitised_con_name %>
ExecStart=/usr/local/bin/manage-user-pod.rb start <%= @pod_yaml_path %> <%= @pod_system_yaml_path %>
ExecStop=-/usr/bin/podman pod stop <%= @sanitised_con_name %>
ExecStop=-/usr/bin/podman pod rm -f <%= @sanitised_con_name %>

<% if @configuration['systemd'] -%>
<%= @configuration['systemd'].map{|k,v| "#{k}=#{v}" }.join("\n") %>

<% end -%>
# security settings
PrivateTmp=true
<% if @real_homedir !~ /^\/home\// -%>
ProtectHome=true
InaccessibleDirectories=/home
<% end -%>

<% if @facts['os']['release']['major'].to_i >= 8 -%>
ProtectSystem=strict
<% else -%>
ProtectSystem=full
<% end -%>

ReadOnlyDirectories=/
ReadWriteDirectories=<%= @real_homedir %>
ReadWriteDirectories=/var/lib/containers/users/<%= @user %>/run/
ReadWriteDirectories=/var/lib/containers/users/<%= @user %>/data/
ReadWriteDirectories=/var/lib/containers/users/<%= @user %>/storage/
ReadWriteDirectories=/run/pods/<%= @uid %>/
ReadWriteDirectories=/proc
<% if @facts['os']['release']['major'].to_i >= 8 -%>
ReadWriteDirectories=/dev/shm
<% end -%>
<% @real_volumes.each do |host_dir,container_dir|
  if container_dir =~ /.*:(Z,)?r(w|o)(,Z)?$/ -%>
ReadWriteDirectories=<%= host_dir %>
<% end
end -%>

[Install]
WantedBy=multi-user.target
