# templates/user-container.service.erb
[Unit]
Description=Container <%= @name %>
After=network.target

[Service]
Type=simple
KillMode=none
Restart=on-failure

User=<%= @user %>
Group=<%= @real_group %>

Environment="XDG_RUNTIME_DIR=/run/pods/<%= @uid %>"

SyslogIdentifier=<%= @unique_name %>

ExecStart=/var/lib/containers/users/<%= @user %>/bin/<%= @unique_name %>.sh start
ExecStop=/var/lib/containers/users/<%= @user %>/bin/<%= @unique_name %>.sh stop

# security settings
PrivateTmp=true
<% unless @real_homedir !~ /^home\// -%>
ProtectHome=true
InaccessibleDirectories=/home
<% end -%>

ProtectSystem=full
# to be migrated once in EL7
# ProtectSystem=strict

# not yet available
# to be migrated once in EL7
#PrivateUsers=true
#ProtectKernelTunables=true
#ProtectKernelModules=true
#ProtectControlGroups=yes

ReadOnlyDirectories=/
ReadWriteDirectories=<%= @real_homedir %>
ReadWriteDirectories=/var/lib/containers/users/<%= @user %>/run/
ReadWriteDirectories=/var/lib/containers/users/<%= @user %>/storage/
ReadWriteDirectories=/run/pods/<%= @uid %>/
<% @real_volumes.each do |host_dir,container_dir|
  if container_dir =~ /.*:rw(,Z)?$/ -%>
ReadWriteDirectories=<%= host_dir %>
<% end
end -%>

[Install]
WantedBy=multi-user.target
