# templates/user-pod.service.erb
[Unit]
Description=<%= @name %> Podman Container
After=network.target

[Service]
Type=simple
KillMode=none
RemainAfterExit=yes
Restart=on-failure

User=<%= @user %>
Group=<%= @real_group %>

Environment="XDG_RUNTIME_DIR=/run/pods/<%= @uid %>"

SyslogIdentifier=<%= @unique_name %>

ExecStartPre=-/usr/bin/podman pod rm -f <%= @sanitised_con_name %>
<% if @run_flags['skip_image_pull'] -%>
ExecStartPre=-/usr/local/bin/pod-update-image.sh <%= @pod_yaml_path %>
<% end -%>
ExecStart=/usr/bin/podman play kube <%= @pod_yaml_path %>

ExecStop=-/usr/bin/podman pod stop <%= @sanitised_con_name %>
ExecStop=-/usr/bin/podman pod rm -f <%= @sanitised_con_name %>

ExecReload=-/usr/bin/podman pod stop <%= @sanitised_con_name %>
ExecReload=-/usr/bin/podman pod rm -f <%= @sanitised_con_name %>

<% if @configuration['systemd'] -%>
<%= @configuration['systemd'].map{|k,v| "#{k}=#{v}" }.join("\n") %>

<% end -%>
[Install]
WantedBy=multi-user.target
