# templates/user-api-socket.service.erb
[Unit]
Description=Podman API Service for <%= @name %>
After=network.target

[Service]
Delegate=true
Type=exec
KillMode=process
Restart=always

User=<%= @name %>
Group=<%= @real_group %>

<% if @manage_xdg_runtime -%>
Environment="XDG_RUNTIME_DIR=/run/pods/<%= @uid %>"
<% end -%>
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/<%= @uid %>/bus"

SyslogIdentifier=<%= @unique_name %>

ExecStart=/usr/bin/podman system service --timeout=0

# security settings
PrivateTmp=true
<% if @real_homedir !~ /^\/home\// -%>
ProtectHome=true
InaccessibleDirectories=/home
<% end -%>

ProtectSystem=strict

ReadOnlyDirectories=/
ReadWriteDirectories=<%= @real_homedir %>
ReadWriteDirectories=/var/lib/containers/users/<%= @name %>/run/
ReadWriteDirectories=/var/lib/containers/users/<%= @name %>/data/
ReadWriteDirectories=/var/lib/containers/users/<%= @name %>/storage/
ReadWriteDirectories=/var/lib/containers/users/<%= @name %>/tmpdir/
<% Array(@readwrite_dirs).each do |rw| -%>
ReadWriteDirectories=<%= rw %>
<% end -%>
ReadWriteDirectories=/run/pods/<%= @uid %>/
<% unless @manage_xdg_runtime -%>
ReadWriteDirectories=/tmp
<% end -%>
ReadWriteDirectories=/proc
ReadWriteDirectories=/dev/shm

[Install]
WantedBy=multi-user.target
