# templates/user-api-socket.service.erb
[Unit]
Description=Podman API Service for <%= @name %>
After=network.target

[Service]
Delegate=true
Type=exec
KillMode=process
Restart=always

User=<%= @name %>
Group=<%= @real_group %>

Environment="XDG_RUNTIME_DIR=/run/pods/<%= @uid %>"

SyslogIdentifier=<%= @unique_name %>

ExecStart=/usr/bin/podman system service --timeout=0

# security settings
PrivateTmp=true
<% if @real_homedir !~ /^\/home\// -%>
ProtectHome=true
InaccessibleDirectories=/home
<% end -%>

ProtectSystem=strict

ReadOnlyDirectories=/
ReadWriteDirectories=<%= @real_homedir %>
ReadWriteDirectories=/var/lib/containers/users/<%= @name %>/run/
ReadWriteDirectories=/var/lib/containers/users/<%= @name %>/data/
ReadWriteDirectories=/var/lib/containers/users/<%= @name %>/storage/
ReadWriteDirectories=/var/lib/containers/users/<%= @name %>/tmpdir/
ReadWriteDirectories=/run/pods/<%= @uid %>/
ReadWriteDirectories=/proc
ReadWriteDirectories=/dev/shm

[Install]
WantedBy=multi-user.target
