# templates/user-container.service.erb
[Unit]
Description=Pod <%= @name %>

[Service]
KillMode=none

User=<%= @user %>
Group=<%= @group %>

Environment="XDG_RUNTIME_DIR=/run/pods/<%= @uid %>"

SyslogIdentifier=<%= @unique_name %>

ExecStart=/var/lib/containers/users/<%= @user %>/bin/<%= @sanitised_title %>.sh start
ExecStop=/var/lib/containers/users/<%= @user %>/bin/<%= @sanitised_title %>.sh stop

[Install]
WantedBy=multi-user.target
