# templates/cron/cron.timer.epp
<%- | String $on_calendar = 'daily',
      String $randomized_delay_sec = '1d',
| -%>
[Timer]
OnCalendar=<%= $on_calendar %>
<% if $randomized_delay_sec { -%>
RandomizedDelaySec=<%= $randomized_delay_sec %>
<% } -%>

[Install]
WantedBy=timers.target
